# EventBus/selenium-grid-eventbus.conf
; Documentation of this file format -> http://supervisord.org/configuration.html

[program:selenium-grid-eventbus]

priority=0
command=/opt/bin/start-selenium-grid-eventbus.sh
autostart=true
autorestart=false
startsecs=0
startretries=0

;Logs (all Hub activity redirected to stdout so it can be seen through "docker logs"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
